# debian/postinst
#!/bin/sh
set -e

systemctl enable connman.service

#DEBHELPER#

exit 0
